# scripts/crontab
0 3 * * * /scripts/backup_pg.sh >> /var/log/backup.log 2>&1
